# supervisord.conf
[supervisord]
nodaemon=true
user=hindsight
logfile=/dev/stdout
logfile_maxbytes=0

[program:nginx]
command=/usr/local/bin/logpipe.sh nginx 12 1 "" /usr/local/bin/start-nginx.sh
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stopasgroup=true
killasgroup=true

[program:hindsight-api]
command=/usr/local/bin/logpipe.sh hindsight-api 60 5 "WORKER_STATS" hindsight-api
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stopasgroup=true
killasgroup=true

[program:control-plane]
command=/usr/local/bin/logpipe.sh control-plane 12 1 "" hindsight-control-plane
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stopasgroup=true
killasgroup=true
environment=PORT="9999",HOSTNAME="0.0.0.0",HINDSIGHT_CP_DATAPLANE_API_URL="http://127.0.0.1:8889",HINDSIGHT_CP_DATAPLANE_API_KEY="%(ENV_HINDSIGHT_API_TENANT_API_KEY)s"

[program:ops-proxy]
# Standalone admin service: GET /logs (durable-log retrieval). Admin-key gated
# (LOGS_ADMIN_KEY), rate-limited, read-only. Runs via system python3.11.
command=/usr/local/bin/logpipe.sh ops-proxy 12 1 "" /usr/local/bin/python3.11 -m uvicorn app:app --host 127.0.0.1 --port 8001
directory=/opt/ops-proxy
user=hindsight
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
startsecs=5
startretries=3
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=LOGS_ADMIN_KEY="%(ENV_LOGS_ADMIN_KEY)s",LOGS_DIR="/data/logs"
